# services/cpld-update@.service

[Unit]
Description=Update CPLD %I

[Service]
Type=oneshot
RemainAfterExit=no
Environment="ARGS=%I"
ExecStart=cpldupdate $ARGS
