# services/orin-update@.service

[Unit]
Description=Update Orin

[Service]
Type=oneshot
RemainAfterExit=no
Environment="ARGS=%I"
ExecStart=/usr/bin/orin-flash.sh $ARGS
