# services/psu-update@.service

[Unit]
Description=Update PSU %I

[Service]
Type=oneshot
RemainAfterExit=no
Environment="ARGS=%I"
ExecStart=psufwupgrade fwupgrade $ARGS
