# services/rt-update@.service

[Unit]
Description=Update RT %I

[Service]
Type=oneshot
RemainAfterExit=no
Environment="ARGS=%I"
ExecStart=/bin/echo To update $ARGS
