# services/rt-update-direct@.service

[Unit]
Description=Update RT %I

[Service]
Type=oneshot
RemainAfterExit=no
Environment="ARGS=%I"
ExecStart=/usr/bin/aries-update $ARGS
