# services/smcu-update@.service

[Unit]
Description=Update SMCU

[Service]
Type=oneshot
RemainAfterExit=no
Environment="ARGS=%I"
ExecStart=/usr/bin/smcu-flash.sh $ARGS
